# hatest/rootfs/etc/services.d/hatest/finish
#!/usr/bin/execlineb -S0
# Niks nodig; placeholder voor nette shutdown
exit 0
