# hatest/rootfs/etc/services.d/hatest/run
#!/usr/bin/with-contenv bash
set -euo pipefail
mkdir -p /run/apache2 /var/log/apache2 \
         /var/www/localhost/htdocs \
         /var/www/localhost/cgi-bin
exec /usr/sbin/httpd -D FOREGROUND


