# hatest/rootfs/etc/services.d/hatest/run
#!/usr/bin/with-contenv bash
set -euo pipefail
exec /usr/local/bin/httpd-foreground
