# hatest/rootfs/etc/services.d/hatest/run
#!/usr/bin/with-contenv bash
set -euo pipefail

# Vereiste runtime-dirs voor Alpine's apache2
mkdir -p /run/apache2 /var/log/apache2 /var/www/localhost/htdocs /var/www/cgi-bin

# Start Apache in foreground (zodat s6 hem kan beheren)
exec /usr/sbin/httpd -D FOREGROUND

